# modules/named/files/etc/systemd/system/master-rsync.timer
[Unit]
Description=rsync master DNS data

[Timer]
OnUnitActiveSec=1hours
Persistent=True

[Install]
WantedBy=basic.target
