# modules/proxy/files/etc/systemd/system/update-squidGuard.timer
[Unit]
Description=Download and install updated squidGuard blacklist

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
