# modules/baseconfig/files/etc/systemd/system/freshclam.timer
[Unit]
Description=Freshclam

[Timer]
OnCalendar=00/1:45
Persistent=True

[Install]
WantedBy=basic.target