# modules/baseconfig/files/etc/systemd/system/puppet-agent.service
[Unit]
Description=Run Puppet Agent to apply configuration changes

[Service]
Type=oneshot
ExecStart=/usr/bin/env puppet agent
