# modules/packages/files/etc/systemd/system/aide.service
[Unit]
Description=AIDE Security Configuration Check

[Service]
Type=oneshot
ExecStart=/usr/bin/env aide --check
