# modules/packages/files/etc/systemd/system/centos-base-repo-update.timer
[Unit]
Description=YUM weekly update

[Timer]
OnCalendar=Sun *-*-* 00:30:00
Persistent=true

[Install]
WantedBy=basic.target
