# modules/packages/files/etc/systemd/system/freshclam.service
[Unit]
Description=Freshclam

[Service]
Type=oneshot
ExecStart=/usr/bin/freshclam --quiet
