# modules/packages/files/etc/systemd/system/rkhunter-update.timer
[Unit]
Description=RKHunter update

[Timer]
OnCalendar=00/1:45
Persistent=True

[Install]
WantedBy=basic.target
