# modules/packages/files/etc/systemd/system/rkhunter-update.service
[Unit]
Description=RKHunter update

[Service]
Type=oneshot
ExecStart=/bin/rkhunter --propupd
