# modules/packages/files/etc/systemd/system/yum-daily-update.timer
[Unit]
Description=YUM daily update

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=True

[Install]
WantedBy=basic.target
